# release/axis-backend/axis-backend.service
[Unit]
Description=Axis TrendAI Backend
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=axis
Group=axis
WorkingDirectory=/opt/axis/backend
EnvironmentFile=/opt/axis/backend/.env
ExecStart=/opt/axis/backend/.venv/bin/python -m uvicorn app.main:app --host 127.0.0.1 --port 8000 --proxy-headers
Restart=always
RestartSec=5
NoNewPrivileges=true

[Install]
WantedBy=multi-user.target
